# crates/neovex-bin/src/machine/assets/neovex.socket.tmpl
[Unit]
Description=Neovex API Socket

[Socket]
ListenStream={guest_neovex_socket}
SocketMode=0600
DirectoryMode=0755

[Install]
WantedBy=sockets.target
